# systemd/knxd.socket
[Unit]
Description=KNX Daemon (socket)

[Socket]
ListenStream=/var/run/knxd
ListenStream=6720

[Install]
WantedBy=sockets.target
